# systemd/.config/systemd/user/dyn-wallpaper.service
[Unit]
Description=Update Dynamic Wallpaper

[Service]
ExecStart=/home/bazyli/bin/dyn-wallpaper

[Install]
WantedBy=default.target
